# appliance/dispatch-firstboot.service
[Unit]
Description=Dispatch appliance first-boot setup (PostgreSQL + per-VM secrets)
# Run once; the script writes the marker when it succeeds.
ConditionPathExists=!/var/lib/dispatch/.appliance-firstboot-done
# Before the relay so PostgreSQL is configured and reachable before Dispatch's first start.
Before=dispatch.service
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/dispatch-appliance/firstboot.sh
# Surface first-boot output on the serial console (and journal) so CI/operators can see what happened.
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
