# installer/linux/dispatch.service
[Unit]
Description=Dispatch SMTP Relay
After=network.target

[Service]
# notify: the service signals readiness via sd_notify (Program.cs UseSystemd) so systemd knows when it's up.
Type=notify
User=dispatch
Group=dispatch
# Binaries live in /opt/dispatch; appsettings.json + the spool live under the data dir (the content root).
WorkingDirectory=/var/lib/dispatch
ExecStart=/opt/dispatch/Dispatch.Service
Environment=DOTNET_ENVIRONMENT=Production
Environment=DISPATCH_LOG_DIR=/var/log/dispatch
Restart=always
RestartSec=5
# Bind the privileged SMTP ports (25, 587) as the unprivileged 'dispatch' user. Without this the
# listener can't bind them and falls back to 2525. AmbientCapabilities grants the capability to the
# process; CapabilityBoundingSet limits it to just this one (and it must be in the bounding set for the
# ambient grant to take). NoNewPrivileges stays true — ambient caps are granted by systemd at exec and
# are unaffected by it.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
# Hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
# Spool .eml/.meta files contain message content (spec §14.5): force mode 600 (rw-------) on every
# file the service creates so message bodies are never group/world-readable.
UMask=0177
ReadWritePaths=/var/lib/dispatch /var/log/dispatch

[Install]
WantedBy=multi-user.target
