# installer/linux/dispatch.service
[Unit]
Description=Dispatch SMTP Relay
After=network.target

[Service]
Type=simple
User=dispatch
Group=dispatch
# Binaries live in /opt/dispatch; appsettings.json + the spool live under the data dir (the content root).
WorkingDirectory=/var/lib/dispatch
ExecStart=/opt/dispatch/Dispatch.Service
Environment=DOTNET_ENVIRONMENT=Production
Environment=DISPATCH_LOG_DIR=/var/log/dispatch
Restart=always
RestartSec=5
# Hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
# Spool .eml/.meta files contain message content (spec §14.5): force mode 600 (rw-------) on every
# file the service creates so message bodies are never group/world-readable.
UMask=0177
ReadWritePaths=/var/lib/dispatch /var/log/dispatch

[Install]
WantedBy=multi-user.target
